# openvpn/fs/usr/lib/systemd/system/openvpn-client-monitor@.path
[Path]
Unit=openvpn-client-monitor@.service
PathChanged=/etc/openvpn/client/client-out.conf

[Install]
WantedBy=multi-user.target
